# airootfs/etc/skel/.config/autostart/karch_theme.desktop
[Desktop Entry]
Name=KArch Theme
Exec=karch_theme
icon=utilities-terminal
Type=Application
Categories=System;Utility;
