# work/x86_64/airootfs/etc/xinetd.d/tftp
service tftp
{
	socket_type	= dgram
	protocol	= udp
	wait		= yes
	user		= nobody
	server		= /usr/bin/tftpd
	server_args	= /usr/var/tftpboot
	disable		= yes
}
