# work/x86_64/airootfs/var/lib/pacman/local/iputils-20200821-1/files
%FILES%
etc/
etc/xinetd.d/
etc/xinetd.d/tftp
usr/
usr/bin/
usr/bin/arping
usr/bin/clockdiff
usr/bin/ninfod
usr/bin/ping
usr/bin/rarpd
usr/bin/rdisc
usr/bin/tftpd
usr/bin/tracepath
usr/lib/
usr/lib/systemd/
usr/lib/systemd/system/
usr/lib/systemd/system/ninfod.service
usr/lib/systemd/system/rarpd@.service
usr/lib/systemd/system/rdisc.service
usr/share/
usr/share/locale/
usr/share/locale/de/
usr/share/locale/de/LC_MESSAGES/
usr/share/locale/de/LC_MESSAGES/iputils.mo
usr/share/locale/fr/
usr/share/locale/fr/LC_MESSAGES/
usr/share/locale/fr/LC_MESSAGES/iputils.mo
usr/share/locale/ja/
usr/share/locale/ja/LC_MESSAGES/
usr/share/locale/ja/LC_MESSAGES/iputils.mo
usr/share/locale/pt_BR/
usr/share/locale/pt_BR/LC_MESSAGES/
usr/share/locale/pt_BR/LC_MESSAGES/iputils.mo
usr/share/locale/tr/
usr/share/locale/tr/LC_MESSAGES/
usr/share/locale/tr/LC_MESSAGES/iputils.mo
usr/share/locale/uk/
usr/share/locale/uk/LC_MESSAGES/
usr/share/locale/uk/LC_MESSAGES/iputils.mo
usr/share/locale/zh_CN/
usr/share/locale/zh_CN/LC_MESSAGES/
usr/share/locale/zh_CN/LC_MESSAGES/iputils.mo
usr/share/man/
usr/share/man/man8/
usr/share/man/man8/arping.8.gz
usr/share/man/man8/clockdiff.8.gz
usr/share/man/man8/ninfod.8.gz
usr/share/man/man8/ping.8.gz
usr/share/man/man8/rarpd.8.gz
usr/share/man/man8/rdisc.8.gz
usr/share/man/man8/tftpd.8.gz
usr/share/man/man8/tracepath.8.gz

