# airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
#
# SPDX-License-Identifier: GPL-3.0-or-later

[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin bhole --noclear %I 38400 linux
